# systemd/pummelchen-release-health.timer
[Unit]
Description=Run Pummelchen release health monitor every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
