# systemd/pummelchen-tested-updates.service
[Unit]
Description=Pummelchen Tested Updates Feed Generator
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/var/minecraft_mods
ExecStart=/usr/bin/flock -n /var/lock/pummelchen-tested-updates.lock /usr/bin/python3 /var/minecraft_mods/scripts/tested_updates_worker.py --db /var/minecraft_mods/data/minecraft_mods.sqlite --output /var/minecraft_mods/site/public/tested-updates.json --days 30
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=full
ReadWritePaths=/var/minecraft_mods/site/public /var/minecraft_mods/data/minecraft_mods.sqlite
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target