# systemd/pummelchen-acceptance-lab-cleanup.timer
[Unit]
Description=Daily cleanup of mod acceptance lab artifacts
Requires=pummelchen-acceptance-lab-cleanup.service

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target