# systemd/pummelchen-release-health.service
[Unit]
Description=Pummelchen release health monitor
After=network-online.target pummelchen-minecraft.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/var/minecraft_mods
ExecStart=/usr/bin/python3 /var/minecraft_mods/scripts/release_health_monitor.py --db /var/minecraft_mods/data/minecraft_mods.sqlite --server-dir /var/minecraft_26.1.2 --release-root /var/minecraft_mods/releases --public-downloads /var/minecraft_mods/site/public/downloads --base-url http://127.0.0.1:7788 --service pummelchen-minecraft.service --quiet
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/var/minecraft_mods/site/public /var/minecraft_mods/site/public/downloads /var/minecraft_26.1.2/client-package/mods
