# systemd/pummelchen-minecraft.service
[Unit]
Description=Pummelchen Minecraft server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/minecraft_26.1.2
ExecStart=/var/minecraft_26.1.2/start.sh
Restart=on-failure
RestartSec=20
TimeoutStopSec=180
SuccessExitStatus=0 143
LimitNOFILE=65535
UMask=0022
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=full
ReadWritePaths=/var/minecraft_26.1.2 /var/minecraft_mods

[Install]
WantedBy=multi-user.target
