# systemd/pummelchen-client-log-receiver.service
[Unit]
Description=Pummelchen client log upload receiver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/minecraft_mods
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/python3 /var/minecraft_mods/scripts/client_log_receiver.py --db /var/minecraft_mods/data/minecraft_mods.sqlite --upload-dir /var/minecraft_mods/client_log_uploads --token-file /var/minecraft_mods/secrets/client-log-upload.token --host 127.0.0.1 --port 7791 --max-mb 25
Restart=always
RestartSec=5
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/minecraft_mods/data /var/minecraft_mods/client_log_uploads /var/minecraft_mods/secrets
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
