# systemd/pummelchen-live-stats.service
[Unit]
Description=Pummelchen live stats JSON writer

[Service]
Type=oneshot
WorkingDirectory=/var/minecraft_mods
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/python3 /var/minecraft_mods/scripts/live_stats_feed.py --output /var/minecraft_mods/site/public/live-stats.json --state /var/minecraft_mods/site/live-stats-history.json --server-dir /var/minecraft_26.1.2
UMask=0022
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/minecraft_mods/site
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
