# systemd/pummelchen-acceptance-lab-cleanup.service
[Unit]
Description=Pummelchen Mod Acceptance Lab Cleanup
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/var/minecraft_mods
ExecStart=/usr/bin/python3 /var/minecraft_mods/scripts/cleanup_acceptance_lab.py --lab-root /var/minecraft_mods/mod_acceptance_lab --max-age-days 7
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=full
ReadWritePaths=/var/minecraft_mods/mod_acceptance_lab
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7